# goc/etc/cron.d/cvmfs
0,15,30,45 * * * * root test ! -f /var/lib/oasis/installed || cvmfs_server snapshot -ai
5,20,35,50 * * * * root test ! -f /var/lib/oasis/installed || (PATH=$PATH:/usr/sbin:/usr/share/oasis; generate_replicas; replicate_whitelists) >>/var/log/cvmfs/generate_replicas.log 2>&1
*/3 * * * * root /usr/share/oasis/oasis_replica_status >/dev/null 2>&1
#8 5 * * * root /usr/share/oasis/gc-all-collectable >>/var/log/cvmfs/gc.log 2>&1
0 9 * * * root find /srv/cvmfs/*.*/data/txn -name "*.*" -mtime +2 2>/dev/null|xargs rm -f
