# goc/etc/cron.d/oasis
* * * * *   root   test ! -f /var/lib/oasis/installed || /usr/share/oasis/do_oasis_update
*/3 * * * * root   /usr/share/oasis/oasis_status_stamp >/dev/null 2>&1

