# goc/etc/cron.d/oasis
* * * * *   root   test ! -f /var/lib/oasis/installed || /usr/share/oasis/do_oasis_update
4,19,34,49 * * * * root test ! -f /var/lib/oasis/installed || ! cvmfs_server masterkeycard -k >/dev/null 2>&1 || (PATH=$PATH:/usr/sbin:/usr/share/oasis; generate_whitelists) >>/var/log/oasis/generate_whitelists.log 2>&1
*/3 * * * * root   /usr/share/oasis/oasis_status_stamp >/dev/null 2>&1
0 5 * * *   root   (date;resign_osg_whitelist) >>/var/log/oasis/resign_osg_whitelist.log 2>&1
